# packaging/contacts-service.service
[Unit]
Description=Start the Contacts service helper service

[Service]
Type=simple
ExecStart=/usr/bin/contacts-service-ipcd

[Install]
WantedBy=default.target
