# packaging/contacts-service.socket
[Unit]
Description=Contacts Service socket

[Socket]
ListenStream=/tmp/.contacts_svc_ipc
SmackLabelIPIn=contacts-service
SmackLabelIPOut=contacts-service
Service=contacts-service.service

[Install]
WantedBy=sockets.target
