# packaging/contacts-service.socket
[Unit]
Description=Contacts Service socket

[Socket]
ListenStream=/tmp/.contacts_svc_ipc-monitor

[Install]
WantedBy=sockets.target
