# packaging/contacts-service.service
[Unit]
Description=Start the Contacts service helper service

[Service]
Environment="CMD=timeout" "TIME=90"
ExecStart=/usr/bin/contacts-service-ipcd $CMD $TIME

